# vault/files/upstart.conf
description "Vault server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  exec /usr/local/bin/vault server\
    -config="/etc/vault/config.hcl"
end script
